# duckdns-deb/DEBIAN/postinst
#!/bin/bash
set -e

duckdns-config

# Habilita o serviço e inicia automaticamente
systemctl daemon-reload
systemctl enable duckdns.timer
systemctl start duckdns.timer
