# duckdns-deb/DEBIAN/postrm
#!/bin/bash
set -e

# Desativa o serviço e remove arquivos extras
systemctl stop duckdns.timer
systemctl disable duckdns.timer
systemctl daemon-reexec
